# deploy/systemd/injection-monitor.service
[Unit]
Description=Injection molding camera monitor (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/injection-monitor/backend
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/pi/injection-monitor/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
